# skillshub-web/deploy/skillshub-web.service
[Unit]
Description=SkillsHub Web (Next.js)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/skillshub/skillshub-web
EnvironmentFile=/opt/skillshub/skillshub-web/.env.production
ExecStart=/usr/bin/npm run start -- -p 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
